# misc/image-loader.timer
[Unit]
Description=Run the image-loader every 5 minutes

[Timer]
OnCalendar=*:0/5
